# files/supervisor_consul.conf
[program:consul]
directory=/tmp/consul
command=/usr/local/bin/consul agent -config-dir /etc/consul.d 
numprocs=1
autostart=true
autorestart=true
stopsignal=INT
